# optional/cpulimit/cpulimit.service
[Unit]
Description=cpulimit
After=network-online.target

[Service]
ExecStart=/bin/cpulimit.sh

[Install]
WantedBy=multi-user.target
